# packages/crystalfontz/crystalfontz-5.0.0.tar.gz/crystalfontz-5.0.0/systemd/crystalfontz.service
[Unit]
Description=Crystalfontz LCD
Requires=dbus.socket
After=dbus.socket

[Service]
Type=dbus
BusName=org.jfhbrook.crystalfontz
ExecStart=/usr/bin/python3 -m crystalfontz.dbus.service --system
Restart=on-failure

[Install]
WantedBy=multi-user.target
